# code/environments/production/modules/batch/templates/batch.cron.erb
#!/usr/bin/env zsh
# This deferjob is managed by puppet, do not modify.

PATH="/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin"

<%= @command %>

